# network/netperf/programs.conf
[program:iperf3]
command=/usr/bin/iperf3 -sV
user=perf
